# config/ghost-pool.service
[Unit]
Description=Ghost Pool - Decentralized Bitcoin Mining Pool
Documentation=https://github.com/bitcoin-ghost/ghost-pool
After=network-online.target bitcoind.service
Wants=network-online.target

[Service]
Type=simple
User=ghost
Group=ghost

# Main executable
ExecStart=/opt/ghost/bin/ghost-pool --config /etc/ghost/pool.toml --data-dir /home/ghost/.ghost

# Restart configuration
# - always: restart on any exit (including clean shutdown)
# - on-failure: restart only on non-zero exit codes
Restart=on-failure
RestartSec=5

# Exit code 100 triggers forced restart (used by config update API)
# When ghost-pool exits with 100, systemd treats it as success but still restarts
RestartForceExitStatus=100
SuccessExitStatus=100

# Environment
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Allow write access to data directory
ReadWritePaths=/home/ghost/.ghost
ReadWritePaths=/var/log/ghost

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Timeouts
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
