# config/ghostd.service
[Unit]
Description=Ghost Core - Bitcoin Ghost Full Node
Documentation=https://github.com/bitcoin-ghost
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ghost
Group=ghost

EnvironmentFile=-/etc/ghost/ghostd.env
ExecStart=/opt/ghost/bin/ghostd \
  -signet \
  -datadir=/var/lib/bitcoin \
  -server=1 \
  -rpcbind=127.0.0.1 \
  -rpcport=38332 \
  -rpcuser=${GHOST_RPC_USER} \
  -rpcpassword=${GHOST_RPC_PASSWORD} \
  -rpcallowip=127.0.0.1/32 \
  -maxconnections=32 \
  -listen=1 \
  -port=38333 \
  -ghostreaper=enabled \
  -zmqpubhashblock=tcp://127.0.0.1:28332 \
  -zmqpubhashtx=tcp://127.0.0.1:28333 \
  -zmqpubsequence=tcp://127.0.0.1:28334

# Restart configuration
Restart=on-failure
RestartSec=30

# Resource limits
LimitNOFILE=65536

# Security hardening
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
