# scripts/checkpoint/ghost-checkpoint.service
[Unit]
Description=Ghost checkpoint generation

[Service]
Type=oneshot
User=ghost
ExecStart=/opt/ghost/bin/generate-checkpoint.sh
